# deploy/5217.service
[Unit]
Description=5217 Flask App
After=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/5217
EnvironmentFile=-/etc/5217/5217.env
ExecStart=/home/ubuntu/5217/.venv/bin/python -m gunicorn --bind 0.0.0.0:5000 app:app
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
